# deploy/docker/sysrepo-netopeer2/latest/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=debug

[program:netopeer2-server]
command=/usr/bin/netopeer2-server -d
autorestart=true
redirect_stderr=true
priority=3
