# deploy/docker/end-to-end/dockerfiles/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=debug

[program:sysrepod]
command=/usr/local/bin/sysrepod -d
autorestart=true
redirect_stderr=true
priority=1

[program:sysrepo-plugind]
command=/usr/local/bin/sysrepo-plugind -d
autorestart=true
redirect_stderr=true
priority=2

[program:netopeer2-server]
command=/usr/local/bin/netopeer2-server -d
autorestart=true
redirect_stderr=true
priority=3

[program:ietf-interfaces]
command=/opt/dev/sysrepo/build/examples/application_example ietf-interfaces
autorestart=true
redirect_stderr=true
priority=4
